# smartcar.service
[Unit]
Description=Smart Car IoT Service
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/smart-car
Environment=ADA_USERNAME=your_username_here ADA_KEY=your_key_here
ExecStart=/usr/bin/python3 /home/pi/smart-car/src/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
